# crates/command-executor/tests/systemd-container/test-services/echo-service/etc/systemd/system/echo-service.service
[Unit]
Description=Test Echo Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/echo-service.sh
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
